# sample.cron
#!/bin/bash
ruby /etc/mysql_backup/mysql_backup.rb backup my_database > /dev/null
exit 0